# deploy/yolo-detect.service
[Unit]
Description=Orchard Patrol Vision desktop launcher
After=graphical.target network-online.target
Wants=network-online.target
ConditionPathExists=/home/rm/orchard-patrol-vision/detect.sh

[Service]
Type=simple
User=rm
Group=rm
WorkingDirectory=/home/rm/orchard-patrol-vision
Environment=DISPLAY=:0
Environment=QT_QPA_PLATFORM=xcb
Environment=DATA_MODE=debug
Environment=TELEMETRY_SERIAL_PORT=/dev/ttyGPS_IN
Environment=TELEMETRY_SERIAL_BAUDRATE=9600
ExecStart=/bin/bash /home/rm/orchard-patrol-vision/detect.sh --foreground
KillMode=control-group
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
