# deploy/yolo-detect.service
[Unit]
Description=Orchard Patrol Vision desktop application
After=graphical-session.target
PartOf=graphical-session.target
ConditionPathExists=/home/rm/orchard-patrol-vision/detect.sh
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=/home/rm/orchard-patrol-vision
Environment=DISPLAY=:0
Environment=QT_QPA_PLATFORM=xcb
Environment=APP_PRESET=client_b
Environment=DATA_MODE=debug
Environment=TELEMETRY_SERIAL_PORT=/dev/ttyGPS_IN
Environment=TELEMETRY_SERIAL_BAUDRATE=9600
EnvironmentFile=-/home/rm/orchard-patrol-vision/.env
ExecStart=/bin/bash /home/rm/orchard-patrol-vision/detect.sh --foreground
KillMode=control-group
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=graphical-session.target
